# etc/systemd/system/snapshot.service
[Unit]
Description=Captura de pantalla local por RTSP

[Service]
Type=oneshot
ExecStart=/usr/local/bin/snapshot.sh

