# etc/systemd/system/snapshot.timer
[Unit]
Description=Timer para ejecutar la caputra

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=10s

[Install]
WantedBy=timers.target
